# conf/cron/prep.cron
10 0 0 * * /home/dafafafa/libexec/cron.pl -t prep
